# jkdc/setup.service
[Unit]
After=mysql.service

[Service]
ExecStart=/root/app-configuration-script.sh

[Install]
WantedBy=default.target
